# deploy/streamhub-backup.timer
[Unit]
Description=Run the StreamHub backup daily
Documentation=file:///opt/streamhub/streamhub-docs/operations/BACKUPS.md

[Timer]
# Daily at 03:15 UTC, +/- up to 5 min jitter. Persistent=true runs a missed
# backup on next boot (e.g. the box was off overnight).
OnCalendar=*-*-* 03:15:00
RandomizedDelaySec=300
Persistent=true
Unit=streamhub-backup.service

[Install]
WantedBy=timers.target
