# deploy/streamhub-core.service
[Unit]
Description=StreamHub Core (NestJS) - LiveKit management API
After=network-online.target redis-server.service livekit.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/streamhub-core
EnvironmentFile=/opt/streamhub-core/.env
ExecStart=/usr/bin/node dist/main.js
Restart=on-failure
RestartSec=3
User=streamhub
LimitNOFILE=200000

[Install]
WantedBy=multi-user.target
